# thickness-monitor.service
[Unit]
Description=Thickness Monitor Merged Backend Service
After=network.target

[Service]
Type=simple
User=linux
Group=linux
WorkingDirectory=/home/linux/final_webapp/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/linux/final_webapp/backend/.venv/bin/python /home/linux/final_webapp/backend/merged_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target