# thickness-monitor.service
[Unit]
Description=Merged Thickness Monitor Backend (KVM Cloud Server)
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/merged/backend
Environment=PYTHONUNBUFFERED=1
Environment=CLOUD_MODE=true
Environment=SERVER_PORT=5002
Environment=INGEST_API_KEY=merged-secret-2026
Environment=DB_HOST=localhost
Environment=DB_NAME=sensor_db
Environment=DB_USER=rapl
Environment=DB_PASS=rapl2026
ExecStart=/opt/merged/venv/bin/python merged_server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
