# agent/thickness-agent.service
[Unit]
Description=Thickness Agent (CD22 sensor reader + cloud uploader)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Headless config (optional): DEVICE_ID / DEVICE_KEY / SENSORS / SERVER_URL live
# in this file. The leading "-" means "don't fail if it's missing" (wizard path).
EnvironmentFile=-/etc/thickness-agent/agent.env
Environment=THICKNESS_AGENT_CONFIG=/etc/thickness-agent/config.json
ExecStart=/opt/thickness-agent/thickness-agent
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
